# crontab
*/2 * * * * /tmp/copyfeed.sh >> /rl/data/logs/facebook-auto-feed/cron.log 2>&1
